# supervisor.conf
[program:your_program]
command=python main.py
autostart=true
autorestart=true
stdout_logfile=log/app.log
stderr_logfile=log/app.err
